# node/systemd/ticker-node-bundled@.service
[Unit]
Description=ticker.cash node — bundled notary+publisher, slot %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=TICKER_HOME=%h/.ticker-slot-%i
Environment=TICKER_ELECTRUM_ACCEPT_ANY_CERT=1
ExecStart=%h/.local/bin/ticker-node --notary --publisher
Restart=on-failure
RestartSec=15s
Slice=ticker.slice
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
