# scripts/systemd/coordinator/ticker.slice
[Unit]
Description=ticker.cash service slice (CPU-capped)
Before=slices.target

[Slice]
CPUQuota=180%
