# scripts/systemd/coordinator/ticker-node-bundled@.service
[Unit]
Description=ticker.cash coordinator — bundled notary+publisher, slot %i
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/toorik/projects/ticker.cash/daemon
Environment=PATH=/home/toorik/.nvm/versions/node/v24.12.0/bin:/usr/bin:/bin
Environment=TICKER_ELECTRUM_HOST=159.195.80.214
Environment=TICKER_ELECTRUM_PORT=50001
ExecStart=/home/toorik/projects/ticker.cash/daemon/node_modules/.bin/tsx scripts/ticker-node.ts --notary --publisher --slot %i
Restart=on-failure
RestartSec=15s
Slice=ticker.slice

[Install]
WantedBy=default.target
